# tests/ci/rc.local
#!/bin/bash

cd /opt/pyroute2

[ "$1" = "test" ] && {
    for p in python2 python3; do {
        which $p && make test skip=test_stress python=$p 2>&1 | tee test-$p.log
        pushd tests
            for f in *xml; do {
                mv $f $p-$f
            } done
        popd
    } done
    sync
    poweroff
} || {
    nohup $0 test &
}
